# ci-tools/fpga-image/gha-runner.service

[Unit]
After=serial-getty.target
After=network-online.target
Wants=network-online.target
Conflicts=serial-getty@ttyPS0.service

[Service]
Type=oneshot
StandardInput=tty-force
StandardOutput=tty-force
StandardError=tty-force
TTYPath=/dev/ttyPS0
User=runner
WorkingDirectory=/home/runner
ExecStart=/home/runner/launch-runner.sh

[Install]
WantedBy=multi-user.target
